# rootfs/etc/services.d/ublox-gps/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup the u-blox GPS RTK service
# ==============================================================================

bashio::log.info "Stopping u-blox GPS RTK service..."
